# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10
startinsecs=2

[program:reverb]
; Servidor WebSocket de Laravel Reverb.
; Escucha en el puerto interno 8081 (nginx proxea /app → ws://127.0.0.1:8081).
; Se inicia DESPUÉS de nginx/php-fpm para que la app ya esté lista.
command=php /var/www/html/artisan reverb:start --host=0.0.0.0 --port=8081 --no-interaction
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
startinsecs=5
